# scripts/systemd/agent-deck-creds-refresh.service
# agent-deck keep-warm OAuth refresh daemon (systemd --user unit).
#
# Part B of the recurring "Please run /login" outage fix. Keeps each profile's
# canonical ~/.claude*/.credentials.json access token fresh so the N symlinked
# worker sessions never hit an expired token and never race Anthropic's
# single-use rotating refresh token. Subscription-safe: it uses the profile's
# existing OAuth refresh token + clientId — NO API key, no per-token billing.
#
# ─────────────────────────────────────────────────────────────────────────────
# NOT auto-started. Enable it deliberately, once you have confirmed both
# profiles are logged in at their canonical path:
#
#   # 1. (one time) make sure canonical is logged in:
#   CLAUDE_CONFIG_DIR=~/.claude      claude   # then /login
#   CLAUDE_CONFIG_DIR=~/.claude-work claude   # then /login
#
#   # 2. install + enable the unit:
#   mkdir -p ~/.config/systemd/user
#   cp scripts/systemd/agent-deck-creds-refresh.service ~/.config/systemd/user/
#   #   (edit ExecStart below if `agent-deck` is not on the default PATH)
#   systemctl --user daemon-reload
#   systemctl --user enable --now agent-deck-creds-refresh.service
#
#   # 3. watch it:
#   systemctl --user status agent-deck-creds-refresh.service
#   journalctl --user -u agent-deck-creds-refresh.service -f
#
# To keep it running after logout (so it warms tokens while you are away):
#   loginctl enable-linger "$USER"
#
# By default it keeps ~/.claude and ~/.claude-work warm (whichever has a
# .credentials.json). To pin specific profiles, add repeated --config-dir flags
# to ExecStart, e.g.:
#   ExecStart=%h/.local/bin/agent-deck creds-refresh --config-dir %h/.claude
# ─────────────────────────────────────────────────────────────────────────────

[Unit]
Description=agent-deck keep-warm OAuth credential refresh (subscription-safe)
Documentation=https://github.com/asheshgoplani/agent-deck
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Adjust the path if `agent-deck` lives elsewhere (e.g. /usr/local/bin or a Go
# bin dir). `%h` expands to the user's home directory.
ExecStart=%h/.local/bin/agent-deck creds-refresh --interval 25m --threshold 20m
Restart=always
RestartSec=30
# The token files are 0600 secrets; keep the process unprivileged and sandboxed.
NoNewPrivileges=true
ProtectSystem=strict
PrivateTmp=true
# It must WRITE the canonical credentials under the profile dirs AND create the
# proper-lockfile lock dir Claude uses, which is a SIBLING of each profile dir
# in $HOME (realpath(CONFIG_DIR)+".lock", e.g. ~/.claude.lock). Because that
# lock is created at runtime directly in $HOME, the whole home must be writable;
# ProtectSystem=strict still keeps the rest of the filesystem read-only.
ReadWritePaths=%h

[Install]
WantedBy=default.target
